# crontab
# DO NOT EDIT THIS FILE ON THE SERVER. THIS FILE IS MANAGEMENTED UNDER GIT.
PATH=/home/dotcloud/perl5/bin:/opt/perl5/perls/current/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/usr/games
PERL5LIB=/home/dotcloud/perl5/lib/perl5

# m h  dom mon dow   command

*/5 *  *   *   *     semaphore --max=1 --file=/tmp/crawler.semaphore -- cronlog perl /home/dotcloud/code/script/cron/crawler.pl
