# crontab
# DO NOT EDIT THIS FILE ON THE SERVER. THIS FILE IS MANAGEMENTED UNDER GIT.

# m h  dom mon dow   command

*/5 *  *   *   *     semaphore --max=1 --file=/tmp/crawler.semaphore -- cronlog perl /home/dotcloud/code/script/cron/crawler.pl
